# files/etc/systemd/system/system76-cameraoff.service
[Unit]
Description=System76 Camera Off

[Service]
ExecStart=/usr/bin/system76-cameraoff

[Install]
WantedBy=multi-user.target
